# deploy/aio/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:postgres]
command=su -c "postgres -D /var/lib/postgresql/data -k /tmp" postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/postgres.log
stderr_logfile=/var/log/postgres.log

[program:redis]
command=redis-server --appendonly no
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/redis.log
stderr_logfile=/var/log/redis.log

[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/caddy.log
stderr_logfile=/var/log/caddy.log

[program:uvicorn]
command=uv run uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2
directory=/app
autostart=true
autorestart=true
priority=40
startsecs=5
stdout_logfile=/var/log/uvicorn.log
stderr_logfile=/var/log/uvicorn.log

[program:worker-fast]
command=uv run arq app.worker.fast.WorkerSettings
directory=/app
autostart=true
autorestart=true
priority=50
stdout_logfile=/var/log/worker-fast.log
stderr_logfile=/var/log/worker-fast.log

[program:worker-slow]
command=uv run arq app.worker.slow.WorkerSettings
directory=/app
autostart=true
autorestart=true
priority=50
stdout_logfile=/var/log/worker-slow.log
stderr_logfile=/var/log/worker-slow.log
